# configs/autostart/firefox.desktop
[Desktop Entry]
Type=Application
Name=Firefox
Icon=firefox
Exec=firefox %u
